# resources/mkosi/signstar/mkosi.extra/usr/lib/systemd/system/systemd-sysupdate.service.d/override.conf
[Service]
ExecStart=
# NOTE: OpenPGP verification currently disabled until https://github.com/systemd/mkosi/commit/bcffc33d9c268c7e33dbc479f4dd04eae0437949 is released
# Also: GnuPG is bizarre and *requires* the keyring file to *not* be ASCII armored, because reasons 🙄
ExecStart=/usr/lib/systemd/systemd-sysupdate update --reboot --verify=no
