# resources/mkosi/signstar/mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/systemd-repart.service.d/override.conf
[Unit]
After=sysroot.mount
ConditionDirectoryNotEmpty=|/sysroot/usr/lib/repart.d
